# qdb-server/control/postinst
#!/bin/sh

USER="qdb"
GROUP="$USER"
ETC_PATH=/etc/qdb
SHARE_PATH=/usr/share/qdb
DB_PATH=/var/lib/qdb
LOG_PATH=/var/log/qdb

# Database
mkdir -p $DB_PATH
chown $USER:$GROUP $DB_PATH

# Log
mkdir -p $LOG_PATH
touch $LOG_PATH/qdbd.log
chown -R $USER:$GROUP $LOG_PATH

# share
mkdir -p $SHARE_PATH
chown -R $USER:$GROUP $SHARE_PATH

# configuration dir
mkdir -p $ETC_PATH

# key generation

if [ ! -f $SHARE_PATH/cluster_public.key ]
then
    echo "Performing new install, initializing cluster keys"

    qdb_cluster_keygen -s $ETC_PATH/cluster_private.key -p $SHARE_PATH/cluster_public.key
    chown $USER:$GROUP $ETC_PATH/cluster_private.key
    chmod 0600 $ETC_PATH/cluster_private.key
    chown $USER:$GROUP $SHARE_PATH/cluster_public.key
    chmod 0644 $SHARE_PATH/cluster_public.key
fi


if [ ! -f $ETC_PATH/qdbsh_private.key ]
then
    # add user for qdbsh
    echo "Creating qdbsh user"
    qdb_user_add -u qdbsh -p $ETC_PATH/users.conf -s $ETC_PATH/qdbsh_private.key

    chmod 0644 $ETC_PATH/users.conf
    chown $USER:$GROUP $ETC_PATH/users.conf

    chmod 0600 $ETC_PATH/qdbsh_private.key
    chown $USER:$GROUP $ETC_PATH/qdbsh_private.key
fi

if [ ! -f $ETC_PATH/qdbd.conf ]
then
    # Config
    echo "Generating default configuration.."
    qdbd --gen-config \
         --log-directory=$LOG_PATH \
         --storage-engine=rocksdb \
         --rocksdb-root=$DB_PATH  \
         --security=true \
         --cluster-private-file=$ETC_PATH/cluster_private.key \
         --user-list=$ETC_PATH/users.conf \
         > $ETC_PATH/qdbd.conf

    chown $USER:$GROUP $ETC_PATH/qdbd.conf
fi

echo "Ensuring file permissions are correct"
chown $USER:$GROUP $ETC_PATH/qdbsh_private.key
chmod 0600 $ETC_PATH/qdbsh_private.key

echo "Setting system control variables"
sysctl -e -q -p /etc/sysctl.d/30-quasardb.conf 2>&1 || true

echo "Setting up systemd"

mkdir -p /usr/lib/systemd/system
cp -f /usr/share/qdb/systemd/qdbd.service /usr/lib/systemd/system/qdbd.service

systemctl daemon-reload 2>&1 || true
systemctl enable qdbd.service
systemctl start qdbd.service 2>&1 || true
